# .config/systemd/user/btwatch.path
[Unit]
Description=Automatically add torrents from watch directory
BindsTo=btpd.service
After=btpd.service

[Path]
PathExistsGlob=%h/downloads/*.torrent
Unit=btwatch.service

[Install]
WantedBy=paths.target
